# tests/crontab
# CDC Validation Cron Schedule
# Format: minute hour day month weekday command

# Run validation every hour at minute 0 (source env vars first)
1 * * * * export $(cat /app/.env_for_cron | xargs) && cd /app && /usr/local/bin/python3 tests/validation.py --all >> /proc/1/fd/1 2>&1
